# sponsorblockcast/rootfs/etc/services.d/sponsorblockcast/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
bashio::log.info 'Starting sponsorblockcast...'

# Set sponsorblockcast env params
SBCPOLLINTERVAL=$(bashio::config 'poll_interval')
SBCSCANINTERVAL=$(bashio::config 'scan_interval')
SBCCATEGORIES=$(bashio::config 'categories')
SBCYOUTUBEAPIKEY=$(bashio::config 'youtube_api_key')
SBCDEBUG=$(bashio::config 'debug')

# Run sponsorblockcast
exec env SBCPOLLINTERVAL="${SBCPOLLINTERVAL}" SBCSCANINTERVAL="${SBCSCANINTERVAL}" SBCCATEGORIES="${SBCCATEGORIES}" SBCYOUTUBEAPIKEY="${SBCYOUTUBEAPIKEY}" SBCDEBUG="${SBCDEBUG}" /usr/bin/sponsorblockcast.sh