# sponsorblockcast/rootfs/etc/services.d/sponsorblockcast/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2191
# ==============================================================================
# Home Assistant Community Add-on: sponsorblockcast
# Runs the sponsorblockcast
# ==============================================================================
bashio::log.info 'Starting sponsorblockcast...'

# Set sponsorblockcast env params
SBCPOLLINTERVAL=$(bashio::config 'poll_interval')
SBCSCANINTERVAL=$(bashio::config 'scan_interval')
SBCCATEGORIES=$(bashio::config 'categories')
SBCYOUTUBEAPIKEY=$(bashio::config 'youtube_api_key')

# Run sponsorblockcast
exec /usr/bin/sponsorblockcast.sh